# caps2esc.conf
start on startup
task
exec /home/jake/bin/caps2esc.sh
